# vllm.service
# vllm.service -- Systemd user service for vLLM + Qwen3.6-35B-A3B AWQ 4-bit
# DGX Spark (GB10, Grace Blackwell)
#
# Installation (as normal user -- no root):
#   mkdir -p ~/.config/systemd/user/
#   cp vllm.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now vllm
#
# Check status:
#   systemctl --user status vllm
#
# View logs:
#   journalctl --user -u vllm -f
#
# Stop / Restart:
#   systemctl --user stop vllm
#   systemctl --user restart vllm

[Unit]
Description=vLLM Server -- Qwen3.6-35B-A3B (DGX Spark)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/vllm

# Overridable via: systemctl --user edit vllm
Environment=PORT=8000
Environment=HOST=0.0.0.0
Environment=MAX_MODEL_LEN=262144
Environment=GPU_MEM_UTIL=0.65
Environment=LANGUAGE_ONLY=true
Environment=NUM_SPEC_TOKENS=2
Environment=MODEL_REPO=cyankiwi/Qwen3.6-35B-A3B-AWQ-4bit

# Server startup script (sources .venv)
ExecStart=/srv/vllm/vllm-server.sh

# Graceful shutdown -- vLLM needs time to drain
TimeoutStopSec=120
KillSignal=SIGTERM
SendSIGKILL=no

# Auto-restart on crash
Restart=on-failure
RestartSec=15

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=default.target